# docs/examples/learn-ukrainian.codex-inbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>learn-ukrainian.codex-inbox</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/your-user/projects/learn-ukrainian/.venv/bin/python</string>
    <string>/Users/your-user/projects/learn-ukrainian/scripts/ai_agent_bridge/__main__.py</string>
    <string>inbox</string>
    <string>run</string>
    <string>codex</string>
    <string>--until-idle</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Users/your-user/projects/learn-ukrainian/.agent/wake/codex</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
